# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/fpp-firmware.list
/
/etc/
/etc/modules.d/
/etc/init.d/
/etc/init.d/fpp_enable
/etc/init.d/fpp
/lib/
/lib/firmware/
/lib/firmware/fpp.axf
