# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/genexis-gaps.list
/
/etc/
/etc/init.d/
/etc/init.d/gaps
/usr/
/usr/bin/
/usr/bin/gapsd
/usr/bin/gxConfig
