# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/genexis-autoswitching.list
/
/etc/
/etc/init.d/
/etc/init.d/autoswitching
/usr/
/usr/bin/
/usr/bin/autoswitch
/usr/bin/drg700_autoswitch.sh
