# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/iptables.list
/
/etc/
/etc/firewall.user
/etc/config/
/etc/config/nat
/etc/config/nat_new
/etc/config/route
/etc/config/firewall
/etc/config/firewall_new
/etc/init.d/
/etc/init.d/nat
/etc/init.d/route
/etc/init.d/firewall
/lib/
/lib/firewall/
/lib/firewall/configure_firewall
/lib/firewall/configure_nat_rules
/lib/firewall/configure_firewall6
/lib/firewall/configure_filter_rules
/usr/
/usr/lib/
/usr/lib/firewall.awk
/usr/lib/iptables/
/usr/lib/iptables/libxt_mac.so
/usr/lib/iptables/libxt_standard.so
/usr/lib/iptables/libxt_multiport.so
/usr/lib/iptables/libipt_LOG.so
/usr/lib/iptables/libxt_limit.so
/usr/lib/iptables/libxt_tcp.so
/usr/lib/iptables/libxt_TCPMSS.so
/usr/lib/iptables/libxt_udp.so
/usr/lib/iptables/libipt_icmp.so
/usr/lib/iptables/libipt_REJECT.so
/usr/sbin/
/usr/sbin/iptables-multi
/usr/sbin/iptables
