# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/kmod-fci.list
/
/etc/
/etc/init.d/
/etc/init.d/fci
/lib/
/lib/modules/
/lib/modules/2.6.33.5-drgos-hrg1000-1.6.5-R/
/lib/modules/2.6.33.5-drgos-hrg1000-1.6.5-R/fci.ko
/usr/
/usr/lib/
/usr/lib/libfci.so
/usr/lib/libfci.so.0.0.0
/usr/lib/libfci.so.0
