# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/kmod-packetfront-catv.list
/
/etc/
/etc/modules.d/
/etc/modules.d/61-packetfront-catv
/etc/init.d/
/etc/init.d/catv
/lib/
/lib/modules/
/lib/modules/2.6.33.5-drgos-hrg1000-1.6.5-R/
/lib/modules/2.6.33.5-drgos-hrg1000-1.6.5-R/packetfront-catv.ko
