# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/miniupnpd.list
/
/etc/
/etc/config/
/etc/config/upnpd
/etc/init.d/
/etc/init.d/miniupnpd
/etc/hotplug.d/
/etc/hotplug.d/iface/
/etc/hotplug.d/iface/10-miniupnpd
/usr/
/usr/bin/
/usr/bin/iptables_display.sh
/usr/bin/iptables_removeall.sh
/usr/bin/iptables_flush.sh
/usr/bin/miniupnpd
/usr/bin/iptables_init.sh
