# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/packetfront-sysload.list
/
/etc/
/etc/init.d/
/etc/init.d/sysload
/sbin/
/sbin/sysloadd
