# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/packetfront-igmp.list
/
/etc/
/etc/init.d/
/etc/init.d/igmp
/usr/
/usr/bin/
/usr/bin/igmp
