# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/packetfront-sip.list
/
/etc/
/etc/config/
/etc/config/sip
/etc/init.d/
/etc/init.d/sip
/usr/
/usr/local/
/usr/local/etc/
/usr/local/etc/vapi.conf
/usr/bin/
/usr/bin/sipapp
/usr/bin/debug_voip
/usr/bin/sip_mon
/usr/sbin/
