# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/packetfront-oxsh.list
/
/etc/
/etc/profile
/etc/iso3166-countries.txt
/etc/init.d/
/etc/init.d/ox-startup-config
/usr/
/usr/bin/
/usr/bin/wlan-show
/usr/bin/oxsh
