# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/packetfront-cwmp.list
/
/etc/
/etc/tr104_notify
/etc/config/
/etc/config/tr104
/etc/init.d/
/etc/init.d/tr104
/etc/init.d/tr069_client
/etc/pf_cwmp.param
/etc/dimark.pem
/etc/dps.param
/usr/
/usr/bin/
/usr/bin/cwmpclient
