# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/packetfront-maconfig.list
/
/etc/
/etc/init.d/
/etc/init.d/maconfig
/usr/
/usr/bin/
/usr/bin/maconfig
