# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/radvd.list
/
/etc/
/etc/config/
/etc/config/radvd
/etc/init.d/
/etc/init.d/radvd
/etc/hotplug.d/
/etc/hotplug.d/iface/
/etc/hotplug.d/iface/05-radvd
/usr/
/usr/sbin/
/usr/sbin/radvd
