# Titanium-52/drgos/initramfs/usr/lib/ipkg/info/watchdog.list
/
/etc/
/etc/watchdog.conf
/etc/init.d/
/etc/init.d/wdog
/etc/sysconfig/
/etc/sysconfig/watchdog
/usr/
/usr/bin/
/usr/bin/systemcheck.sh
/usr/sbin/
/usr/sbin/watchdog
