# Titanium-52/drgos/initramfs/etc/hotplug.d/iface/00-netstate
case "$ACTION" in                               
    ifup)
        uci_set_state network "$INTERFACE" up 1 
        
	[ "$STATE" != "renew" ] && {
	    # Fix Bug #20614: possible to ping LAN from WAN
            echo "nat.general.toggle=1" > /tmp/mafifo
	    echo "nat.general.toggle=0" > /tmp/mafifo
  	} 
	;;
esac    

