# Titanium-52/drgos/initramfs/etc/hotplug.d/iface/40-dhcp6c
#!/bin/sh

#pieingCurrently hotplug is not being used for IPv6 and DHCP
#DHCP6 handling is done by dhcp6c-state(dynamic config) 
#and network6(static config)

return 0
