# Titanium-52/drgos/initramfs/etc/hotplug.d/iface/40-pf-update
# Copyright 2011 PacketFront International AB

case "$ACTION" in                               
	ifup|update)
		include /lib/network
		scan_interfaces

                #SNMP
                src_intf_snmp=$(uci -q get snmpd.@agent[0].agentinterface)      
                [ -z "$src_intf_snmp" ] && src_intf_snmp=$(uci -q get network.default.wan_if)
                #only restart on src interface
                [ "$src_intf_snmp" = "$INTERFACE" ] && {
                     #the STATE parameter tells the state of the dhcp, i.e. bound, renew, ...
                     #only restart snmp on bound
                     [ "$STATE" = "bound" ] && {
                         #reload snmpd via mgmt framwork
                         echo "snmpd.all.toggle=1" > /tmp/mafifo
                         echo "snmpd.all.toggle=0" > /tmp/mafifo
                     }                 
                }  

		  #IGMP, tell snooper to reload if igmp upstream interface 
		  #has bound to a new ip. When statically assigned the signal
		  #will be sent from the network script.
		  [ "$STATE" = "bound" ] && {
		    src_intf_igmp=$(uci -q get network.lan.igmpupstreamintf)
		    if [ -n "$src_intf_igmp" ] && [ "$src_intf_igmp" = "$INTERFACE" ]; then
		        pid=`pidof igmp`
  		     	[ -n "$pid" ] && kill -SIGHUP $pid
		    fi	   
		  }
		;;                                
esac    

